# hourly-ping.timer
[Unit]
Description=Run ping service every hour

[Timer]
Unit=ping.service
OnCalendar=hourly

[Install]
WantedBy=default.target
