# ping.service
[Unit]
Description=Play ping sound

[Service]
Type=oneshot
WorkingDirectory=/home/elio/dotfiles
ExecStart=play ping.ogg
